# root/etc/update-motd.d/00-retro-home-header
#!/bin/sh

[ -r /etc/os-release ] && . /etc/os-release

printf "Welcome to Retro Home 🕹  \n\nBuilt with %s (%s %s %s)\n" "${PRETTY_NAME}" "$(uname -o)" "$(uname -r)" "$(uname -m)"
